# infrastructure/pi-deployment/systemd/lmrc-launcher.service
[Unit]
Description=LMRC Application Launcher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/lmrc/shared/scripts/launcher.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
